# PKGBUILDS/hyprmo/hyprmo-utils/hyprmo-setpermissions.service
[Unit]
Description=hyprmo-pinephone
[Service]
ExecStart=/usr/bin/hyprmo_setpermissions.sh
[Install]
WantedBy=multi-user.target
